# command-center-api.service
[Unit]
Description=Laura Dashboard API Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/var/www/laura-dashboard/api-server
ExecStart=/home/ubuntu/.nvm/versions/node/v22.22.0/bin/node server.js
Restart=on-failure
RestartSec=5
EnvironmentFile=/var/www/laura-dashboard/api-server/.env

[Install]
WantedBy=multi-user.target
